# pulsar.in
#!/bin/sh

# PROVIDE: pulsar
# REQUIRE: DAEMON NETWORKING
# KEYWORD: shutdown
#
# pulsar_config (path):  Set to /usr/local/etc/pulsar/config.json by default.
#		          Config-file to use.

. /etc/rc.subr

name="pulsar"
rcvar="${name}_enable"

load_rc_config $name

: ${pulsar_enable:="NO"}
: ${pulsar_config:="/usr/local/etc/${name}/config.json"}

pidfile="/var/run/${name}.pid"
procname="/usr/local/sbin/${name}"
command="/usr/sbin/daemon"
command_args="-o /var/log/${name}.log -p ${pidfile} ${procname} -d -c ${pulsar_config}"

run_rc_command "$1"
